# systemd/white88-backend.service
[Unit]
Description=White88 Backend Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/white88
Environment=PATH=/var/www/white88/venv/bin
Environment=DJANGO_SETTINGS_MODULE=white88_backend.settings
Environment=SECRET_KEY=your-secret-key-here
Environment=DEBUG=False
Environment=ALLOWED_HOSTS=localhost,127.0.0.1,your-domain.com
Environment=DATABASE_URL=postgresql://white88_user:white88_password@localhost:5432/white88_db
ExecStart=/var/www/white88/venv/bin/gunicorn --bind 0.0.0.0:8000 --workers 3 --timeout 120 white88_backend.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=white88-backend

[Install]
WantedBy=multi-user.target 